# Task 4/crontab
*/1 * * * * /usr/bin/periodictask/postmessage.sh